# docker-samba4/phpldapadmin/supervisor/conf.d/phpldapadmin.conf
[program:phpldapadmin]
priority=20
command=apache2ctl -DFOREGROUND
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log

